# archlinux/rootfs-overlay/etc/systemd/system/serial-getty@ttyGS0.service.d/override.conf
[Unit]
Requires=r11t-usb-gadget.service
After=r11t-usb-gadget.service

[Service]
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no
