# archlinux/rootfs-overlay/etc/systemd/system/bluetooth.service.d/20-r11t.conf
[Unit]
Wants=r11t-bluetooth-setup.service
After=r11t-bluetooth-setup.service
