# archlinux/rootfs-overlay/etc/systemd/system/sddm.service.d/20-r11t.conf
[Unit]
Requires=r11t-local-hardware.service
After=r11t-local-hardware.service r11t-audio.service
