# archlinux/rootfs-overlay/etc/xdg/autostart/r11t-audio-volume.desktop
[Desktop Entry]
Type=Application
Name=Initialize OPPO R11T audio volume
Exec=/usr/lib/r11t/set-initial-audio-volume
OnlyShowIn=KDE;
NoDisplay=true
