# archlinux/rootfs-overlay/etc/systemd/system/NetworkManager.service.d/20-r11t.conf
[Unit]
Wants=r11t-wifi.service
After=r11t-wifi.service
